# Unix/Servicios/Timers/ReiniciarSistemaTranviaBilbao.timer
[Unit]
Description=Reiniciar Sistema Tranvia de Bilbao 

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target